# build/dnf_data/etc/s6-overlay/s6-rc.d/sg_mysql_proxy/run
#!/command/execlineb -P
foreground { /etc/s6-overlay/scripts/on-start sg_mysql_proxy }
with-contenv
s6-notifyoncheck -T 600000 -n 0 -w 500 -t 10000 -c "/home/template/init/lib/s6-runprobe sg_mysql_proxy"
fdclose 3
cd /home/template/init/mysql_proxy
fdmove -c 2 1
/home/template/init/mysql_proxy/start_mysql_proxy.sh
