# dnf.conf
[program:gate]
command= /bin/bash -c "./GateRestart"
directory=/root
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:stun]
command= /bin/bash -c "./df_stun_r start"
directory=/home/neople/stun
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:monitor]
command= /bin/bash -c "./df_monitor_r mnt_siroco start"
directory=/home/neople/monitor
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:manager]
command= /bin/bash -c "./df_manager_r manager start"
directory=/home/neople/manager
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:relay]
command= /bin/bash -c "./df_relay_r relay_200 start"
directory=/home/neople/relay
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:bridge]
command= /bin/bash -c "./df_bridge_r bridge start"
directory=/home/neople/bridge
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:channel]
command= /bin/bash -c "./df_channel_r channel start"
directory=/home/neople/channel
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:dbmw_guild]
command= /bin/bash -c "./df_dbmw_r dbmw_gld_siroco start"
directory=/home/neople/dbmw_guild
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:dbmw_mnt]
command= /bin/bash -c "./df_dbmw_r dbmw_mnt_siroco start"
directory=/home/neople/dbmw_mnt
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:dbmw_stat]
command= /bin/bash -c "./df_dbmw_r dbmw_stat_siroco start"
directory=/home/neople/dbmw_stat
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:auction]
command= /bin/bash -c "./df_auction_r ./cfg/auction_siroco.cfg start ./df_auction_r"
directory=/home/neople/auction
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:point]
command= /bin/bash -c "./df_point_r ./cfg/point_siroco.cfg start df_point_r"
directory=/home/neople/point
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:guild]
command= /bin/bash -c "./df_guild_r gld_siroco start"
directory=/home/neople/guild
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:statics]
command= /bin/bash -c "./df_statics_r stat_siroco start"
directory=/home/neople/statics
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:coserver]
command= /bin/bash -c "./df_coserver_r coserver start"
directory=/home/neople/coserver
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:community]
command= /bin/bash -c "./df_community_r community start"
directory=/home/neople/community
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:gunnersvr]
command= /bin/bash -c "./gunnersvr -t30 -i1"
directory=/home/neople/secsvr/gunnersvr
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:zergsvr_secagent]
command= /bin/bash -c "./secagent"
directory=/home/neople/secsvr/zergsvr
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:zergsvr]
command= /bin/bash -c "./zergsvr -t30 -i1"
directory=/home/neople/secsvr/zergsvr
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:game_siroco11]
command= /bin/bash -c "./df_game_r siroco11 start"
directory=/home/neople/game
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:game_siroco52]
command= /bin/bash -c "./df_game_r siroco52 start"
directory=/home/neople/game
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[group:dnf]
programs=gate,stun,monitor,manager,relay,bridge,channel,dbmw_guild,dbmw_mnt,dbmw_stat,auction,point,guild,statics,coserver,community,gunnersvr,zergsvr_secagent,zergsvr,game_siroco11,game_siroco52
priority=999
