# build/dnf_data/etc/s6-overlay/s6-rc.d/daily_job/run
#!/command/execlineb -P
with-contenv
cd /home/template/init/daily_job
/home/template/init/daily_job/daily_job.sh
