# build/dnf_data/etc/s6-overlay/s6-rc.d/scheduler/run
#!/command/execlineb -P
foreground { /etc/s6-overlay/scripts/on-start scheduler }
with-contenv
cd /home/template/init/scheduler
fdmove -c 2 1
/home/template/init/scheduler/scheduler.sh
