# build/dnf_data/etc/s6-overlay/s6-rc.d/scheduler/run
#!/command/execlineb -P
with-contenv
cd /home/template/init/scheduler
/home/template/init/scheduler/scheduler.sh
