# build/dnf_data/etc/s6-overlay/s6-rc.d/tailscale/run
#!/command/execlineb -P
foreground { /etc/s6-overlay/scripts/on-start tailscale }
with-contenv
cd /home/template/init/tailscale
fdmove -c 2 1
/home/template/init/tailscale/start_tailscale.sh
