# build/dnf_data/etc/supervisor/conf.d/dnf.conf
[program:main_mysql_proxy]
command=/bin/bash -c "/home/template/init/mysql_proxy/start_main_mysql_proxy.sh"
directory=/home/template/init/mysql_proxy
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/main_mysql_proxy.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB

[program:sg_mysql_proxy]
command=/bin/bash -c "/home/template/init/mysql_proxy/start_mysql_proxy.sh"
directory=/home/template/init/mysql_proxy
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/sg_mysql_proxy.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=main_mysql_proxy

[program:daily_job]
command=/bin/bash -c "/home/template/init/daily_job/daily_job.sh"
directory=/home/template/init/daily_job
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/daily_job.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=sg_mysql_proxy

[program:stun]
command=/bin/bash -c "/data/run/start_stun.sh"
directory=/home/neople/stun
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/stun.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=daily_job

[program:monitor]
command=/bin/bash -c "/data/run/start_monitor.sh"
directory=/home/neople/monitor
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/monitor.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=stun

[program:manager]
command=/bin/bash -c "/data/run/start_manager.sh"
directory=/home/neople/manager
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/manager.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=monitor

[program:relay]
command=/bin/bash -c "/data/run/start_relay.sh"
directory=/home/neople/relay
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/relay.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=manager

[program:bridge]
command=/bin/bash -c "/data/run/start_bridge.sh"
directory=/home/neople/bridge
user=root
autostart=true
autorestart=unexpected
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/bridge.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=relay

[program:dbmw_guild]
command=/bin/bash -c "/data/run/start_dbmw_guild.sh"
directory=/home/neople/dbmw_guild
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/dbmw_guild.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=channel

[program:dbmw_mnt]
command=/bin/bash -c "/data/run/start_dbmw_mnt.sh"
directory=/home/neople/dbmw_mnt
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/dbmw_mnt.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=dbmw_guild

[program:dbmw_stat]
command=/bin/bash -c "/data/run/start_dbmw_stat.sh"
directory=/home/neople/dbmw_stat
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/dbmw_stat.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=dbmw_mnt

[program:auction]
command=/bin/bash -c "/data/run/start_auction.sh"
directory=/home/neople/auction
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/auction.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=dbmw_stat

[program:point]
command=/bin/bash -c "/data/run/start_point.sh"
directory=/home/neople/point
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/point.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=auction

[program:guild]
command=/bin/bash -c "/data/run/start_guild.sh"
directory=/home/neople/guild
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/guild.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=point

[program:statics]
command=/bin/bash -c "/data/run/start_statics.sh"
directory=/home/neople/statics
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/statics.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=guild

[program:coserver]
command=/bin/bash -c "/data/run/start_coserver.sh"
directory=/home/neople/coserver
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/coserver.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=statics

[program:community]
command=/bin/bash -c "/data/run/start_community.sh"
directory=/home/neople/community
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/community.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=coserver

[program:gunnersvr]
command=/bin/bash -c "/data/run/start_gunnersvr.sh"
directory=/home/neople/secsvr/gunnersvr
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/gunnersvr.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=community

[program:zergsvr_secagent]
command=/bin/bash -c "/data/run/start_zergsvr_secagent.sh"
directory=/home/neople/secsvr/zergsvr
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/zergsvr_secagent.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=gunnersvr

[program:zergsvr]
command=/bin/bash -c "/data/run/start_zergsvr.sh"
directory=/home/neople/secsvr/zergsvr
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/zergsvr.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=zergsvr_secagent

[program:monitor_ip]
command=/bin/bash -c "/home/template/init/monitor_ip/monitor_ip.sh"
directory=/home/template/init/monitor_ip
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/monitor_ip.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
depend=channel

[group:dnf]
programs=stun,monitor,manager,relay,bridge,dbmw_guild,dbmw_mnt,dbmw_stat,auction,point,guild,statics,coserver,community,gunnersvr,zergsvr_secagent,zergsvr,monitor_ip,main_mysql_proxy,sg_mysql_proxy,daily_job
priority=999
