# build/dnf_data/etc/supervisor/conf.d/gate.conf
[program:llnut_gate]
command=/bin/bash -c "/home/template/init/lib/barrier-wait llnut_gate /home/template/init/run/start_gate.sh"
directory=/usr/bin
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/llnut_gate.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=GATE_ENVIRONMENT
