# build/dnf_data/etc/supervisor/conf.d/gate.conf.template
[program:tongyi_gate]
command=/bin/bash -c "./DnfGateServer"
directory=/root
user=root
autostart=false
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/tongyi_gate.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
