# build/dnf_data/etc/supervisor/conf.d/monitor_ip.conf
[program:monitor_ip]
command=/bin/bash -c "/home/template/init/monitor_ip/monitor_ip.sh"
directory=/home/template/init/monitor_ip
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/monitor_ip.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB

[group:monitor_ip]
programs=monitor_ip
priority=999
