# build/dnf_data/etc/supervisor/conf.d/netbird.conf
[program:netbird]
command=/bin/bash -c "/home/template/init/netbird/start_netbird.sh"
directory=/home/template/init/netbird
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/netbird/netbird.log
redirect_stderr=true
depend=channel
