# build/dnf_data/etc/supervisor/conf.d/starter.conf
[program:starter]
command=/bin/bash -c "/home/template/init/starter/start_program.sh"
directory=/home/template/init/starter
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
