# build/dnf_data/etc/supervisor/conf.d/starter.conf
[program:starter]
command=/bin/bash -c "/home/template/init/starter/start_program.sh"
directory=/home/template/init/starter
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/starter.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
