# build/dnf_data/etc/supervisor/conf.d/tailscale.conf
[program:tailscaled]
command=/bin/bash -c "/home/template/init/tailscale/start_tailscaled.sh"
directory=/home/template/init/tailscale
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/tailscale/tailscaled.log
redirect_stderr=true
depend=channel

[program:tailscale]
command=/bin/bash -c "/home/template/init/tailscale/start_tailscale.sh"
directory=/home/template/init/tailscale
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/tailscale/tailscale.log
redirect_stderr=true
depend=tailscaled

