# build/dnf_data/etc/s6-overlay/s6-rc.d/netbird/run
#!/command/execlineb -P
foreground { /etc/s6-overlay/scripts/on-start netbird }
with-contenv
cd /home/template/init/netbird
fdmove -c 2 1
/home/template/init/netbird/start_netbird.sh
