# build/dnf_data/etc/s6-overlay/s6-rc.d/netbird/run
#!/command/execlineb -P
with-contenv
cd /home/template/init/netbird
/home/template/init/netbird/start_netbird.sh
