# build/dnf_data/etc/s6-overlay/s6-rc.d/cleanup/up
/command/with-contenv bash /etc/s6-overlay/scripts/init.d/cleanup.sh
