# build/dnf_data/etc/s6-overlay/s6-rc.d/game_template/run
#!/command/execlineb -P
with-contenv
fdmove -c 2 1
s6-notifyoncheck -T 600000 -n 0 -w 500 -t 10000 -c "/home/template/init/lib/s6-runprobe __SVC_NAME__"
fdclose 3
cd /home/neople/game
/data/run/start_game.sh __CHANNEL_NUM__ __PROCESS_SEQ__
