# plugin/by-gate/by-gate.conf
[program:by_gate]
command=/bin/bash -c "if [ ! -d "./by-gate" ]; then tar -zxvf by-gate.tar.gz; fi; cd by-gate;rm install.log;./start-by.sh | tee install.log"
directory=/data/conf.d
user=root
autostart=true
autorestart=false
stopasgroup=true
killasgroup=true
