# build/dnf_data/etc/s6-overlay/s6-rc.d/tailscaled/run
#!/command/execlineb -P
with-contenv
cd /home/template/init/tailscale
/home/template/init/tailscale/start_tailscaled.sh
