# build/dnf_data/etc/s6-overlay/s6-rc.d/manager/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck -T 1200000 -n 0 -w 500 -t 10000 -c "/home/template/init/lib/s6-runprobe manager"
fdclose 3
cd /home/neople/manager
/data/run/start_manager.sh
