# build/dnf_data/etc/s6-overlay/s6-rc.d/mysql/run
#!/command/execlineb -P
with-contenv
fdmove -c 2 1
bash /etc/s6-overlay/scripts/mysql-run.sh
