# build/dnf_data/etc/s6-overlay/s6-rc.d/main_mysql_proxy/run
#!/command/execlineb -P
with-contenv
cd /home/template/init/mysql_proxy
/home/template/init/mysql_proxy/start_main_mysql_proxy.sh
