# plugin/dnf-console/dnf-console.conf
[program:dnf-console]
command=/bin/bash -c "if [ ! -d "./dnf-console" ]; then tar -zxvf dnf-console.tgz; fi; cd dnf-console;./main -i && ./main -x"
directory=/data/conf.d
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
