# build/dnf_data/etc/s6-overlay/s6-rc.d/llnut_gate/run
#!/command/execlineb -P
with-contenv
cd /usr/bin
/home/template/init/run/start_gate.sh
