# build/dnf_data/supervisor/channel.conf
[program:channel]
command=/bin/bash -c "/data/channel/start_channel.sh"
directory=/home/neople/channel
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/channel.log
redirect_stderr=true
depend=bridge