# build/dnf_data/supervisor/channel.conf
[program:channel]
command=/bin/bash -c "/home/neople/channel/start_channel.sh"
directory=/home/neople/channel
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
environment=LD_PRELOAD=/data/dp/libhook.so
depend=bridge

[program:game_siroco11]
command=/bin/bash -c "/home/neople/game/start_siroco.sh 11 3"
directory=/home/neople/game
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
environment=LD_PRELOAD=/data/dp/libhook.so
depend=zergsvr

[program:game_siroco52]
command=/bin/bash -c "/home/neople/game/start_siroco.sh 52 5"
directory=/home/neople/game
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
environment=LD_PRELOAD=/data/dp/libhook.so
depend=game_siroco11

[group:dnf_channel]
programs=channel,game_siroco11,game_siroco52
priority=999
