# build/dnf_data/supervisor/dnf.conf
[program:main_mysql_proxy]
command=/bin/bash -c "/home/template/init/mysql_proxy/start_main_mysql_proxy.sh"
directory=/home/template/init/mysql_proxy
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/main_mysql_proxy.log
redirect_stderr=true

[program:sg_mysql_proxy]
command=/bin/bash -c "/home/template/init/mysql_proxy/start_mysql_proxy.sh"
directory=/home/template/init/mysql_proxy
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/sg_mysql_proxy.log
redirect_stderr=true
depend=main_mysql_proxy

[program:daily_job]
command=/bin/bash -c "/data/daily_job/daily_job.sh"
directory=/data/daily_job/
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/daily_job.log
redirect_stderr=true
depend=sg_mysql_proxy

[program:stun]
command=/bin/bash -c "killall -9 df_stun_r;rm -rf pid/*.pid && ./df_stun_r start && sleep 2 && cat pid/*.pid |xargs -n1 -I{} tail --pid={} -f /dev/null"
directory=/home/neople/stun
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
depend=daily_job

[program:monitor]
command=/bin/bash -c "killall -9 df_monitor_r;rm -rf pid/*.pid && ./df_monitor_r server start && sleep 2 && cat pid/*.pid |xargs -n1 -I{} tail --pid={} -f /dev/null"
directory=/home/neople/monitor
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
environment=LD_PRELOAD=/home/template/init/libhook.so
depend=stun

[program:manager]
command=/bin/bash -c "killall -9 df_manager_r;rm -rf pid/*.pid && ./df_manager_r server start && sleep 2 && cat pid/*.pid |xargs -n1 -I{} tail --pid={} -f /dev/null"
directory=/home/neople/manager
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
environment=LD_PRELOAD=/home/template/init/libhook.so
depend=monitor

[program:relay]
command=/bin/bash -c "killall -9 df_relay_r;rm -rf pid/*.pid && ./df_relay_r relay_200 start && sleep 2 && cat pid/*.pid |xargs -n1 -I{} tail --pid={} -f /dev/null"
directory=/home/neople/relay
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
depend=manager

[program:bridge]
command=/bin/bash -c "/data/channel/start_bridge.sh"
directory=/home/neople/bridge
user=root
autostart=true
autorestart=unexpected
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/start_bridge.log
redirect_stderr=true
depend=relay

[program:dbmw_guild]
command=/bin/bash -c "/home/neople/dbmw_guild/start_dbmw_guild.sh"
directory=/home/neople/dbmw_guild
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
environment=LD_PRELOAD=/home/template/init/libhook.so
depend=channel

[program:dbmw_mnt]
command=/bin/bash -c "/home/neople/dbmw_mnt/start_dbmw_mnt.sh"
directory=/home/neople/dbmw_mnt
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
environment=LD_PRELOAD=/home/template/init/libhook.so
depend=dbmw_guild

[program:dbmw_stat]
command=/bin/bash -c "/home/neople/dbmw_stat/start_dbmw_stat.sh"
directory=/home/neople/dbmw_stat
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
environment=LD_PRELOAD=/home/template/init/libhook.so
depend=dbmw_mnt

[program:auction]
command=/bin/bash -c "killall -9 df_auction_r;rm -rf pid/*.pid && ./df_auction_r ./cfg/server.cfg start ./df_auction_r && sleep 5 && cat pid/*.pid |xargs -n1 -I{} tail --pid={} -f /dev/null"
directory=/home/neople/auction
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
depend=dbmw_stat

[program:point]
command=/bin/bash -c "killall -9 df_point_r;rm -rf pid/*.pid && ./df_point_r ./cfg/server.cfg start df_point_r && sleep 5 && cat pid/*.pid |xargs -n1 -I{} tail --pid={} -f /dev/null"
directory=/home/neople/point
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
depend=auction

[program:guild]
command=/bin/bash -c "killall -9 df_guild_r;rm -rf pid/*.pid && ./df_guild_r server start && sleep 2 && cat pid/*.pid |xargs -n1 -I{} tail --pid={} -f /dev/null"
directory=/home/neople/guild
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
environment=LD_PRELOAD=/home/template/init/libhook.so
depend=point

[program:statics]
command=/bin/bash -c "killall -9 df_statics_r;rm -rf pid/*.pid && ./df_statics_r server start && sleep 2 && cat pid/*.pid |xargs -n1 -I{} tail --pid={} -f /dev/null"
directory=/home/neople/statics
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
environment=LD_PRELOAD=/home/template/init/libhook.so
depend=guild

[program:coserver]
command=/bin/bash -c "killall -9 df_coserver_r;rm -rf pid/*.pid && ./df_coserver_r server start && sleep 2 && cat pid/*.pid |xargs -n1 -I{} tail --pid={} -f /dev/null"
directory=/home/neople/coserver
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
environment=LD_PRELOAD=/home/template/init/libhook.so
depend=statics

[program:community]
command=/bin/bash -c "killall -9 df_community_r;rm -rf pid/*.pid && ./df_community_r server start && sleep 2 && cat pid/*.pid |xargs -n1 -I{} tail --pid={} -f /dev/null"
directory=/home/neople/community
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
depend=coserver

[program:gunnersvr]
command=/bin/bash -c "killall -9 gunnersvr;rm -rf pid/*.pid && ./gunnersvr -t30 -i1 && sleep 2 && cat pid/*.pid |xargs -n1 -I{} tail --pid={} -f /dev/null"
directory=/home/neople/secsvr/gunnersvr
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
depend=community

[program:zergsvr_secagent]
command=/bin/bash -c "killall -9 secagent;rm -rf secagent.pid && sleep 2 && ./secagent"
directory=/home/neople/secsvr/zergsvr
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
depend=gunnersvr

[program:zergsvr]
command=/bin/bash -c "killall -9 zergsvr;rm -rf zergsvr.pid && ./zergsvr -t30 -i1;sleep 5 && cat zergsvr.pid |xargs -n1 -I{} tail --pid={} -f /dev/null"
directory=/home/neople/secsvr/zergsvr
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
depend=zergsvr_secagent

[program:monitor_ip]
command=/bin/bash -c "/home/template/init/monitor_ip/monitor_ip.sh"
directory=/home/template/init/monitor_ip
user=root
autostart=true
autorestart=unexpected
exitcodes=0
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/monitor_ip.log
redirect_stderr=true
depend=channel

[group:dnf]
programs=stun,monitor,manager,relay,bridge,dbmw_guild,dbmw_mnt,dbmw_stat,auction,point,guild,statics,coserver,community,gunnersvr,zergsvr_secagent,zergsvr,monitor_ip,main_mysql_proxy,sg_mysql_proxy,daily_job
priority=999
