# build/dnf_data/supervisor/gate.conf
[program:tongyi_gate]
command=/bin/bash -c "./DnfGateServer"
directory=/root
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/data/log/tongyi_gate.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
