# build/dnf_data/supervisor/gate.conf
[program:tongyi_gate]
command=/bin/bash -c "./DnfGateServer > /data/log/tongyi_gate.log"
directory=/root
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
